# deploy/cloudflared/cloudflared-quick-sync.service
[Unit]
Description=Sync Career Counsel Quick Tunnel URL
After=cloudflared-quick.service
Requires=cloudflared-quick.service

[Service]
Type=oneshot
EnvironmentFile=/etc/default/career-counsel-quick-tunnel-sync
ExecStart=/usr/bin/python3 /usr/local/lib/career-counsel/sync-quick-tunnel-url.py
